# scripts/sfizz/systemd/sfizz.socket
[Socket]
ListenFIFO=%t/sfizz.stdin
Service=sfizz.service
PartOf=sfizz.service
